# Admin/script_usb_key/usbpmount.rules
ACTION=="add", KERNEL=="sd?1", RUN+="/etc/init.d/do_it_on_usb_connection.sh"
